# ssh.conf
[program:FakeSSH]
command=<path-to-FakeSSH>/fake.py
directory=<path-to-FakeSSH>
user=root
autorestart=true
redirect_stderr=true
